# nano.cron
SHELL=/bin/bash
0 * * 11 * cd ~/Projects/nanowrimo/2013 && make; git push
