# nano.cron
SHELL=/bin/bash
0 * * 11 * MAILTO=`git config --get user.email`; cd ~/Projects/nanowrimo/`date +%Y` && make; git push
